# patches/systemvm/etc/init.d/domr_webserver
#!/bin/bash
# chkconfig: 345 98 02
# description: Web server that sends passwords to User VMs
#   This file exists in /etc/init.d/ in the routing domain
# @VERSION@

bash /root/run_domr_webserver&
