# patches/systemvm/debian/config/etc/cloud-nic.rules
SUBSYSTEM=="net" KERNEL=="eth*" RUN+="/opt/cloud/bin/cloud-nic.sh $env{ACTION} %k"

