# cloudbridge/packages/config/debian/cloud-bridge.postinst
#!/bin/sh -e

case "$1" in
    configure)
        if ! id cloud > /dev/null 2>&1 ; then
            adduser --system --home /usr/share/cloud --no-create-home \
                --group --disabled-password --shell /bin/sh cloud
        fi

	for i in /usr/share/cloud/bridge 
	do
		chmod 0775 $i
		chown -R cloud.cloud $i
	done
	
	if [ "$2" = "" ] ; then # no recently configured version, this is a first install
		/usr/sbin/update-rc.d cloud-bridge defaults || true
	fi
	
    ;;
esac

#DEBHELPER#
