# systemvm/debian/etc/gro-off.rules
SUBSYSTEM=="net" KERNEL=="eth*" DRIVERS=="virtio_net" RUN+="/opt/cloud/bin/gro-off.sh %k"
