# debian/cloud-console-proxy.postinst
#!/bin/sh -e

case "$1" in
    configure)

	for i in /var/log/cloud/console-proxy
	do
		chmod 0770 $i
	done
	
	if [ "$2" = "" ] ; then # no recently configured version, this is a first install
		/usr/sbin/update-rc.d cloud-console-proxy defaults || true
	fi
	
    ;;
esac

#DEBHELPER#