# systemvm/debian/etc/cron.daily/clear_cache
#!/bin/bash

# clear memory cache to ultimately reduce swapping

sync && echo 1 > /proc/sys/vm/drop_caches
