# vnet/distro/centos/SYSCONFDIR/rc.d/init.d/cloud-vnetd.in
#!/bin/bash

# chkconfig: 35 99 05
# description: CloudStack Virtual Network Daemon

# WARNING: if this script is changed, then all other initscripts MUST BE changed to match it as well

. /etc/rc.d/init.d/functions

whatami=cloud-vnetd

# set environment variables

SHORTNAME="$whatami"
PIDFILE=@PIDDIR@/"$whatami".pid
LOCKFILE=@LOCKDIR@/"$SHORTNAME"
LOGFILE=@LOCALSTATEDIR@/log/cloud/vnetd.log
PROGNAME="CloudStack Virtual Network Daemon"

unset OPTIONS
[ -r @SYSCONFDIR@/sysconfig/"$SHORTNAME" ] && source @SYSCONFDIR@/sysconfig/"$SHORTNAME"
DAEMONIZE=@BINDIR@/@PACKAGE@-daemonize
PROG=@SBINDIR@/$SHORTNAME

start() {
        echo -n $"Starting $PROGNAME: "
	daemon --check=$SHORTNAME --pidfile=${PIDFILE} "$DAEMONIZE" \
		-n "$SHORTNAME" -p "$PIDFILE" -l "$LOGFILE" "$PROG" $OPTIONS
	RETVAL=$?
	echo
	[ $RETVAL = 0 ] && touch ${LOCKFILE}
	return $RETVAL
}

stop() {
	echo -n $"Stopping $PROGNAME: "
	killproc -p ${PIDFILE} $SHORTNAME # -d 10 $SHORTNAME
	RETVAL=$?
	echo
	[ $RETVAL = 0 ] && rm -f ${LOCKFILE} ${PIDFILE}
}

recreate_vnets() {
   for br in `ip link | grep vnbr | cut -d":" -f2| cut --complement -c1`; 
   do 
     vnetid=0000:0000:0000:0000:0000:0000:0000:`echo $br | cut -c5-` 
     vn vnet-create -b $br $vnetid &> /dev/null; 
   done
}

# See how we were called.
case "$1" in
  start)
	start
	sleep 1
	recreate_vnets
	;;
  stop)
	stop
	;;
  status)
        status -p ${PIDFILE} $SHORTNAME
	RETVAL=$?
	;;
  restart)
	stop
	sleep 2
	start
	recreate_vnets
	;;
  condrestart)
	if status -p ${PIDFILE} $SHORTNAME >&/dev/null; then
		stop
		sleep 2
		start
		recreate_vnets
	fi
	;;
  *)
	echo $"Usage: $whatami {start|stop|restart|condrestart|status|help}"
	RETVAL=3
esac

exit $RETVAL

