# systemvm/debian/etc/systemd/system/cloud-early-config.service
[Unit]
Description=CloudStack post-boot patching service using cmdline
DefaultDependencies=no

Before=sshd.service

Requires=local-fs.target network-pre.target
After=local-fs.target

[Install]
WantedBy=multi-user.target

[Service]
Type=oneshot
ExecStart=/opt/cloud/bin/setup/cloud-early-config
RemainAfterExit=no
TimeoutStartSec=5min
