# patches/systemvm/debian/config/etc/profile.d/cloud.sh
if [ "`id -u`" -eq 0 ]; then
  PATH=${PATH}:/opt/cloud/bin
fi
export PATH
